# nekoray/scripts/nekoray-daemon.service.in
[Unit]
Description=NekoRay Headless V2Ray Client Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=@DAEMON_USER@
Group=@DAEMON_GROUP@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/nekoray-daemon --config @CONFIG_DIR@ --port 8080 --verbose
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=@CONFIG_DIR@
PrivateDevices=true
PrivateTmp=true

# For TUN mode, we need CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE

# Environment
Environment="QT_QPA_PLATFORM=offscreen"

[Install]
WantedBy=multi-user.target